# noia-agent-docker/templates/noia-agent.service.j2
[Unit]
Description=NOIA Agent
After=multi-user.target

[Service]
Type=simple
Restart=always
RestartSec=1
ExecStart=/usr/local/bin/noia-agent

[Install]
WantedBy=multi-user.target

